# services/systemd/pi-setup-fleet.service
# User systemd unit — fleet daemon (Linux).
# Prefer: run ./install.sh from your clone; on Linux it writes this file to
#   ~/.config/systemd/user/pi-setup-fleet.service
# with WorkingDirectory set to your actual repo path.
#
# Manual install: copy to ~/.config/systemd/user/, set WorkingDirectory to this
# repo (same directory as package.json). The daemon loads .env.runtime from
# that directory automatically — do not commit .env.runtime.
#
#   systemctl --user daemon-reload
#   systemctl --user enable --now pi-setup-fleet.service
#
# User services do NOT run at boot until you log in (SSH/GUI), unless you enable
# lingering so systemd keeps a user session for this account at boot:
#   loginctl enable-linger "$USER"
#   # or: sudo loginctl enable-linger pi
# Check: loginctl show-user "$USER" -p Linger   → Linger=yes

[Unit]
Description=pi-setup fleet daemon
After=network.target

[Service]
Type=simple
# Default matches install.sh / README: git clone … pi.dev
WorkingDirectory=%h/pi.dev
ExecStart=/usr/bin/env node scripts/fleet-daemon.mjs
Restart=always
RestartSec=5
Environment=PI_SETUP_DAEMON_PORT=4269

[Install]
WantedBy=default.target
